# templates/etc/supervisord.conf.in
[unix_http_server]
file=${buildout:directory}/var/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=${buildout:directory}/var/log/supervisord.log 
logfile_maxbytes=50MB       
logfile_backups=10          
loglevel=info               
pidfile=${buildout:directory}/var/supervisord.pid 
nodaemon=false             
minfds=1024                
minprocs=200               

[supervisorctl]
serverurl=unix://${buildout:directory}/var/supervisor.sock ; use a unix:// URL  for a unix socket
history_file=~/.sc_history  ; use readline history if available

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:redis]
command = ${buildout:directory}/bin/redis-server
redirect_stderr=true 
priority=5

[program:public]
command = ${buildout:bin-directory}/python -m SimpleHTTPServer 8080
directory = ${buildout:directory}/www/public
redirect_stderr=true

[program:private]
command = ${buildout:bin-directory}/python -m SimpleHTTPServer 9090
directory = ${buildout:directory}/www/private
redirect_stderr=true
